# scripts/com.switchocr.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template LaunchAgent for the Switch OCR server.

  Replace REPLACE_ME_HOME with the absolute path to the home directory of the
  user that should run the server (for example, /Users/yourname). After copying
  this file to ~/Library/LaunchAgents/com.switchocr.server.plist, run:

      launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.switchocr.server.plist
      launchctl kickstart -k gui/$(id -u)/com.switchocr.server
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.switchocr.server</string>

	<key>WorkingDirectory</key>
	<string>REPLACE_ME_HOME/Projects/SwitchOCR</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>

	<key>ProgramArguments</key>
	<array>
		<string>REPLACE_ME_HOME/Projects/SwitchOCR/run_server.sh</string>
	</array>

	<key>KeepAlive</key>
	<true/>
	<key>RunAtLoad</key>
	<true/>

	<key>StandardOutPath</key>
	<string>REPLACE_ME_HOME/Library/Logs/SwitchOCR/server.out.log</string>
	<key>StandardErrorPath</key>
	<string>REPLACE_ME_HOME/Library/Logs/SwitchOCR/server.err.log</string>
</dict>
</plist>
