# scripts/com.switchocr.frame-capture.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Optional LaunchAgent for maintaining server/latest_frame.jpg independently
  from /ocr-latest. Replace REPLACE_ME_HOME before installing.
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.switchocr.frame-capture</string>

	<key>WorkingDirectory</key>
	<string>REPLACE_ME_HOME/Projects/SwitchOCR</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>

	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/ssh</string>
		<string>-i</string>
		<string>REPLACE_ME_HOME/.ssh/switchocr_local_ed25519</string>
		<string>-o</string>
		<string>IdentitiesOnly=yes</string>
		<string>-o</string>
		<string>BatchMode=yes</string>
		<string>-o</string>
		<string>StrictHostKeyChecking=accept-new</string>
		<string>-o</string>
		<string>ServerAliveInterval=30</string>
		<string>-o</string>
		<string>ServerAliveCountMax=3</string>
		<string>localhost</string>
		<string>cd "$HOME/Projects/SwitchOCR" &amp;&amp; exec ./scripts/run_frame_capture.sh</string>
	</array>

	<key>KeepAlive</key>
	<true/>
	<key>RunAtLoad</key>
	<true/>

	<key>StandardOutPath</key>
	<string>REPLACE_ME_HOME/Library/Logs/SwitchOCR/frame-capture.out.log</string>
	<key>StandardErrorPath</key>
	<string>REPLACE_ME_HOME/Library/Logs/SwitchOCR/frame-capture.err.log</string>
</dict>
</plist>
