# scripts/com.switchocr.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template LaunchAgent for the Switch OCR server.

  Replace REPLACE_ME_HOME with the absolute path to the home directory of the
  user that should run the server (for example, /Users/yourname). After copying
  this file to ~/Library/LaunchAgents/com.switchocr.server.plist, run:

      launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.switchocr.server.plist
      launchctl kickstart -k gui/$(id -u)/com.switchocr.server

  Why the ssh-to-localhost wrapper?
  SysDVR-Client needs a proper user session/terminal-like environment to
  complete its TCP bridge handshake on macOS. Running it directly under
  launchd (nohup, plain daemon, etc.) reliably fails. Re-entering through ssh
  to localhost gives it the session it needs and KeepAlive restarts it if it
  ever exits.
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.switchocr.server</string>

	<key>WorkingDirectory</key>
	<string>REPLACE_ME_HOME/Projects/SwitchOCR</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>

	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/ssh</string>
		<string>-i</string>
		<string>REPLACE_ME_HOME/.ssh/switchocr_local_ed25519</string>
		<string>-o</string>
		<string>IdentitiesOnly=yes</string>
		<string>-o</string>
		<string>BatchMode=yes</string>
		<string>-o</string>
		<string>StrictHostKeyChecking=accept-new</string>
		<string>-o</string>
		<string>ServerAliveInterval=30</string>
		<string>-o</string>
		<string>ServerAliveCountMax=3</string>
		<string>localhost</string>
		<string>cd "$HOME/Projects/SwitchOCR" &amp;&amp; exec ./run_server.sh</string>
	</array>

	<key>KeepAlive</key>
	<true/>
	<key>RunAtLoad</key>
	<true/>

	<key>StandardOutPath</key>
	<string>REPLACE_ME_HOME/Library/Logs/SwitchOCR/server.out.log</string>
	<key>StandardErrorPath</key>
	<string>REPLACE_ME_HOME/Library/Logs/SwitchOCR/server.err.log</string>
</dict>
</plist>
